# infra/update-watcher/printops-update-watcher.service
[Unit]
Description=PrintOps scheduled update watcher
After=network.target printops-api.service

[Service]
Type=oneshot
User=itadmin
Group=itadmin
ExecStart=/bin/bash /home/itadmin/printops/infra/update-watcher/update-watcher.sh
